# conf/cron.daily/update-tls
#!/bin/sh
letsencrypt renew --agree-tos
service nginx reload